# SOURCES/trackmonitor-1/etc/cron.d/trackmonitor
* * * * * root /etc/trackmonitor/monitoring.sh

