# SOURCES/trackmonitor-1/trackmonitor.service
[Unit]
Description=trackmonitor - service for track linux courses

[Service]
Type=oneshot
ExecStart=/etc/trackmonitor/monitoring.sh

